# server/cron.d/stratummonitor
# m h  dom mon dow	user	command
@reboot ircbot  supybot /home/ircbot/supybot/supybot.conf.example > /dev/null
