# server/cron.d/stratummonitor
# m h  dom mon dow	user	command
@reboot ircbot  supybot /home/ircbot/supybot/supybot.conf.example > /dev/null
@reboot ircbot  touch /var/run/stratummonitor-netscan
*/5 *  *   *   *	root	nmap -sn -n -PM 192.168.178.0/24 | sed -n 's/^MAC Address: \([0-9A-Fa-f:]\+\) .*$/\1/p' > /var/run/stratummonitor-netscan
